# github-build/lib/systemd/system/github-build.target
[Unit]
Description=build all the things
