# github-build/lib/systemd/system/github-build@.service
[Service]
Type=oneshot
ExecStart=/bin/ls

[Install]
WantedBy=github-build.target
