# certgen/lib/systemd/system/certgen.timer
[Timer]
Unit=certgen.target
OnActiveSec=0
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
