# certgen/lib/systemd/system/certgen@.service
[Service]
ExecStart=/usr/local/bin/cert-refresh -C /usr/local/lib/acme %I

[Install]
WantedBy=certgen.target
